# deploy/systemd/remux-extensions.slice
[Unit]
Description=Remux extension resource slice

[Slice]
CPUWeight=400
CPUAccounting=yes
MemoryAccounting=yes
TasksAccounting=yes
