# deploy/systemd/remux-core.slice
[Unit]
Description=Remux core resource slice

[Slice]
CPUWeight=1000
CPUAccounting=yes
MemoryAccounting=yes
TasksAccounting=yes
