# deploy/systemd/remux.slice
[Unit]
Description=Remux resource root

[Slice]
CPUAccounting=yes
MemoryAccounting=yes
TasksAccounting=yes
